# codex/systemd/rclone-mount-gdrive.service
[Unit]
Description=Rclone mount for Google Drive
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=RCLONE_CONFIG=/int/.runtime/cloud-access/rclone.conf
ExecStartPre=/int/tools/codex/cloud_access.sh ensure-dirs
ExecStart=/int/tools/codex/cloud_access.sh mount gdrive /int/cloud/gdrive
ExecStop=/int/tools/codex/cloud_access.sh umount /int/cloud/gdrive
Restart=on-failure
RestartPreventExitStatus=3
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=default.target
