# codex/systemd/rclone-mount-yadisk.service
[Unit]
Description=Rclone mount for Yandex Disk
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=RCLONE_CONFIG=/git/.runtime/cloud-access/rclone.conf
ExecStartPre=/git/tools/codex/cloud_access.sh ensure-dirs
ExecStart=/git/tools/codex/cloud_access.sh mount yadisk /git/cloud/yadisk
ExecStop=/git/tools/codex/cloud_access.sh umount /git/cloud/yadisk
Restart=on-failure
RestartPreventExitStatus=3
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=default.target
