# openclaw/systemd/openclaw-gateway.service.d/gog-runtime.conf
[Service]
Environment=XDG_CONFIG_HOME=/home/leon/.config
Environment=GOG_KEYRING_BACKEND=file
Environment=GOG_KEYRING_PASSWORD=
